# anacron/anacrontab
~/.anacron/anacrontab
------------------------------------------------------

SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/home/user/bin
HOME=/home/user
LOGNAME=user

# period  delay  job-id  command
1  5   user.cron.daily    nice run-parts /home/jbajic/.anacron/cron.daily

7  10  user.cron.weekly   nice run-parts /home/jbajic/.anacron/cron.weekly

@monthly  5  user.cron.monthly  nice run-parts /home/jbajic/.anacron/cron.monthly
